# etc/systemd/system/sshd.service.d/restart.conf
[Service]
RestartSec=5
OOMPolicy=continue
OOMScoreAdjust=-1000
